# meshtastic-mqtt-translate/rootfs/etc/services.d/mqtt-relay/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC1008
# ==============================================================================
# Start the MQTT relay service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# # Add your code here

# # Declare variables
# declare message

# ## Get the 'message' key from the user config options.
# message=$(bashio::config 'message')

MQTT_BROKER=$(bashio::services mqtt "host")
MQTT_USERNAME=$(bashio::services mqtt "username")
MQTT_PASSWORD=$(bashio::services mqtt "password")

if [[ ${MQTT_BROKER} != "null" && ! -z ${MQTT_BROKER} ]]; then
	## if MQTT_BROKER does not have "tcp://" prefix, add it
	if [[ ${MQTT_BROKER} != tcp://* ]]; then
		MQTT_BROKER="tcp://${MQTT_BROKER}"
	fi
	## if MQTT_BROKER does not have port, add default port 1883
	if ! [[ ${MQTT_BROKER} =~ :[0-9]+$ ]]; then
		MQTT_BROKER="${MQTT_BROKER}:1883"
	fi
	export MQTT_BROKER
fi

if [[ ${MQTT_USERNAME} == "null" || -z ${MQTT_USERNAME} ]]; then
	export MQTT_USERNAME=""
fi

if [[ ${MQTT_PASSWORD} == "null" || -z ${MQTT_PASSWORD} ]]; then
	export MQTT_PASSWORD=""
fi

echo "Starting Meshtastic MQTT Relay with the following settings:"
echo "MQTT Host: ${MQTT_BROKER}"
echo "MQTT User: ${MQTT_USERNAME}"
# Note: Not printing MQTT_PASSWORD for security reasons
echo "MQTT Password: ********"

## Run your program
exec /opt/mqtt-relay/meshtastic-mqtt-relay
