# meshtastic-mqtt-translate/rootfs/etc/services.d/mqtt-relay/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC1008
# ==============================================================================
# Start the MQTT relay service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# # Add your code here

# # Declare variables
# declare message

# ## Get the 'message' key from the user config options.
# message=$(bashio::config 'message')

MQTT_BROKER=$(bashio::services mqtt "host")
MQTT_USERNAME=$(bashio::services mqtt "username")
MQTT_PASSWORD=$(bashio::services mqtt "password")

export MQTT_BROKER MQTT_USERNAME MQTT_PASSWORD

echo "Starting Meshtastic MQTT Relay with the following settings:"
echo "MQTT Host: ${MQTT_BROKER}"
echo "MQTT User: ${MQTT_USERNAME}"
# Note: Not printing MQTT_PASSWORD for security reasons
echo "MQTT Password: ********"

## Run your program
exec /opt/mqtt-relay/meshtastic-mqtt-relay
